# templates/peer-scripts/wg-status-server.service.j2
[Unit]
Description=wgmesh status web server for {{ peer.name }}
After=network.target wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/wg-status-server.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
