# templates/peer-scripts/wg-watchdog-task.xml.j2
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler definition for the Porthole WireGuard watchdog.
  Generated by porthole for peer {{ peer.name }}.

  Import with:
    schtasks /Create /TN "Porthole\wg-watchdog-{{ peer.name }}" /XML wg-watchdog-task.xml /F

  Or via Task Scheduler GUI:
    Action → Import Task → select this file.

  The task runs every 5 minutes as SYSTEM (no login required).
  The PowerShell script must be copied to C:\ProgramData\Porthole\wg-watchdog.ps1
  (or adjust the Arguments path below).
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-01-01T00:00:00</Date>
    <Author>Porthole</Author>
    <Description>WireGuard watchdog for peer {{ peer.name }}. Checks hub reachability and restarts WireGuard if the connection drops.</Description>
    <URI>\Porthole\wg-watchdog-{{ peer.name }}</URI>
  </RegistrationInfo>
  <Triggers>
    <!-- Run at system startup (delayed 60 s to allow WireGuard to initialise) -->
    <BootTrigger>
      <Enabled>true</Enabled>
      <Delay>PT60S</Delay>
    </BootTrigger>
    <!-- Run every 5 minutes indefinitely -->
    <TimeTrigger>
      <Enabled>true</Enabled>
      <StartBoundary>2026-01-01T00:00:00</StartBoundary>
      <Repetition>
        <Interval>PT5M</Interval>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
    </TimeTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <!-- Run as SYSTEM — required to restart the WireGuard service -->
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT2M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>PowerShell.exe</Command>
      <Arguments>-NonInteractive -NoProfile -ExecutionPolicy Bypass -File "C:\ProgramData\Porthole\wg-watchdog.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
